# sv/commander/run
#!/bin/sh

exec 2>&1
exec python /restcommander/play-1.2.4/play start /restcommander/AgentMaster
